# .config/profile.d/10-themes.sh
#!/bin/sh
# ~/.config/profile.d/10-themes.sh
# themes configuration
# <alexforsale@yahoo.com>
