# .config/profile.d/99-angular.sh
#!/bin/sh
# ~/.config/profile.d/99-angular.sh
# angular configuration
# <alexforsale@yahoo.com>
command -v ng && source <(ng completion script)
